# fs/ml605_skeleton/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
null::sysinit:/bin/mount -t proc proc /proc
null::sysinit:/bin/mount -o remount,rw / # REMOUNT_ROOTFS_RW
null::sysinit:/bin/mkdir -p /dev/pts
null::sysinit:/bin/mkdir -p /dev/shm
null::sysinit:/bin/mount -a
null::sysinit:/bin/hostname -F /etc/hostname
null::sysinit:/bin/mount -t devpts devpts /dev/pts -o noexec,nosuid
null::sysinit:/bin/mount -t debugfs debugfs /sys/kernel/debug
# now run any rc scripts
::sysinit:/etc/init.d/rcS

# Put a getty on the serial port
ttyUL0::respawn:/sbin/getty -L ttyUL0 57600 vt100 # GENERIC_SERIAL
null::sysinit:mknod /dev/spidev0.2 c 153 0
null::sysinit:mknod /dev/spidev1.2 c 153 1
null::sysinit:bitflip /dev/spidev0.2 /home/default/AD9548\ script.txt
null::sysinit:bitflip /dev/spidev1.2 /home/default/AD9548\ script.txt

# ndso junk
null::sysinit:/sbin/udhcpc -n -t 6
null::sysinit:/usr/sbin/inetd
null::sysinit:/home/default/testxcomm.sh

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
null::shutdown:/etc/init.d/rcK
null::shutdown:/bin/umount -a -r
null::shutdown:/sbin/swapoff -a

